# deb-build/build/DEBIAN/postinst
#!/bin/bash

systemctl daemon-reload

echo RESTARTING SERVICES
if [ "$(systemctl is-enabled journalctl2psql)" == "enabled" ]; then
	echo "starting journalctl2psql"
	systemctl restart journalctl2psql
fi

if [ "$(systemctl is-enabled journalctl2psql-remote)" == "enabled" ]; then
	echo "starting journalctl2psql-remote"
	systemctl restart journalctl2psql-remote
fi

if [ "$(systemctl is-enabled nxlog-udp-json2psql)" == "enabled" ]; then
	echo "starting nxlog-udp-json2psql"
	systemctl restart nxlog-udp-json2psql
fi